# ansible-playbooks/roles/bitbucket/templates/bitbucket.service.j2
[Unit]
Description=Atlassian Bitbucket Server Service
After=syslog.target network.target

[Service]
Type=forking
User={{ bitbucket_user }}
ExecStart={{ bitbucket_dir }}/bin/start-bitbucket.sh
ExecStop={{ bitbucket_dir }}/bin/stop-bitbucket.sh

[Install]
WantedBy=multi-user.target
